# ops/document-simu/vs14-writer-build.timer
[Unit]
Description=Vacation Station 14 — writer rebuild timer
Documentation=file:///opt/vacation-station/docs/OPERATIONS.md

[Timer]
# Sun 04:45 UTC weekly (templates don't change often)
OnCalendar=Sun *-*-* 04:45:00 UTC
RandomizedDelaySec=15min
Persistent=true
Unit=vs14-writer-build.service

[Install]
WantedBy=timers.target
